# ~/.config/systemd/user/virtdev-tmux-status@.service
[Unit]
Description=virtdev tmux status listener for pane %i

[Service]
Type=simple
ExecStart=socat UNIX-LISTEN:%t/virtdev/socket/status-%i.sock,fork,max-children=2,unlink-early "EXEC:%h/.local/bin/virtdev-tmux-status %i"

# Sandbox: the listener processes untrusted data from guest VMs.
# jq (C parser) is the primary attack surface — a crafted JSON
# payload exploiting a jq vulnerability would execute in this
# context. These directives limit what a compromised handler can do.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
MemoryMax=128M
TasksMax=8
