# ~/.config/systemd/user/virtdev-tmux-status@.service
[Unit]
Description=virtdev tmux status listener for pane %i

[Service]
Type=simple
ExecStart=socat UNIX-LISTEN:%t/virtdev/socket/status-%i.sock,fork,max-children=2,unlink-early,unlink-close "EXEC:%h/.local/bin/virtdev-tmux-status %i"

# Sandbox: the listener processes untrusted data from guest VMs.
# Parsing uses bash builtins (read) and sanitization uses tr/sed.
# These directives limit what a compromised handler can do.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
MemoryMax=128M
TasksMax=8
